# applesmc-fan-control.service
[Unit]
Description=Fan control

[Service]
ExecStart=/path/to/applesmc-fan-control

[Install]
WantedBy=multi-user.target
